# docker/config/crontabs/root
* * * * * cd /var/www/html/ && php artisan schedule:run >> /myscript.log 2>&1
